# files/etc/service/fpm-error-forwarder/run
#!/bin/sh
exec tail -F -n 0 /var/log/fpm-php.www.log
