# files/etc/service/fpm-error-forwarder/run
#!/bin/sh
touch /var/log/fpm-php.www.log
exec tail -F -n 0 /var/log/fpm-php.www.log
