# files/etc/service/log-tailer/run
#!/bin/sh
# Tail all relevant logs
touch /var/log/fpm-php.www.log
touch /tmp/blackfire.log
touch /var/log/messages
touch /var/log/error
touch /var/log/debug
touch /var/log/syslog
exec tail -F -n 0 /var/log/fpm-php.www.log /tmp/blackfire.log /var/log/messages /var/log/error /var/log/debug /var/log/syslog
